# rootfs/etc/s6-overlay/s6-rc.d/krb5kdc/run
#!/usr/bin/with-contenv bash
set -e
exec /usr/sbin/krb5kdc -n
