# rootfs/etc/s6-overlay/s6-rc.d/namenode/finish
#!/usr/bin/with-contenv bash
set -eo pipefail
echo "[INFO] namenode stopping (exit: ${1}, signal: ${2})"
exit 0
